# core/systemd/battery-core.service
[Unit]
Description=Battery Core — Anthropic usage monitor background service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/share/battery/node_modules/.bin/battery-core --loop
Restart=on-failure
RestartSec=10
Environment=HOME=%h

[Install]
WantedBy=default.target
